# etc/ot-recorder.service
# /etc/systemd/system/ot-recorder.service

[Unit]
Description=OwnTracks Recorder

[Service]
Type=simple
User=owntracks
WorkingDirectory=/
ExecStart=/usr/sbin/ot-recorder 'owntracks/+/+'

[Install]
WantedBy=multi-user.target
